# etc/systemd/system/babeld.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/sbin/babeld -I '' wlan1 br-lan
