# etc/systemd/system/mesh-web.service
[Unit]
Description=Natak Mesh Web Interface
After=network.target mesh-start.service babeld.service
Wants=babeld.service

[Service]
Type=simple
User=natak
Group=natak
WorkingDirectory=/opt/nucleus/web
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/bin/python3 /opt/nucleus/web/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
