# etc/systemd/system/reticulum-mesh-rns-monitor.service
[Unit]
Description=Reticulum Mesh RNS Monitor
After=reticulum-mesh-reticulum.service
Requires=reticulum-mesh-reticulum.service

[Service]
Type=simple
User=root
Group=root
SupplementaryGroups=dialout
WorkingDirectory=/home/natak/reticulum_mesh/rns_stats
Environment=PYTHONPATH=/home/natak/reticulum_mesh
Environment=HOME=/home/natak
ExecStart=/bin/bash -c 'sleep 10 && /usr/bin/python3 -u rns_monitor.py'
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
